# rootfs/target/FILOGIC880-BE19000/etc/uci-defaults/opensync_mtk_zero_wait_dfs_disable
#!/bin/sh

/usr/bin/datconf -f /etc/wireless/mediatek/mt7990.b1.dat set DfsDedicatedZeroWait "0"
/usr/bin/datconf -f /etc/wireless/mediatek/mt7990.b1.dat set DfsZeroWaitDefault "0"
